# scripts/launchd/com.quanttrading.strategy-allocator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quanttrading.strategy-allocator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/varadbandekar/Documents/Quant_trading/venv/bin/python</string>
        <string>/Users/varadbandekar/Documents/Quant_trading/scripts/strategy_allocator.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/varadbandekar/Documents/Quant_trading</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/varadbandekar/Library/Logs/quanttrading/strategy-allocator.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/varadbandekar/Library/Logs/quanttrading/strategy-allocator.err.log</string>
</dict>
</plist>
